# System/OS - Linux/06 서비스와 소프트웨어/02 서비스 관리/lab02/system-monitor.service
[Unit]
Description=System Monitor Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/system-monitor.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
